# deploy/systemd/tikpal-kiosk.service
[Unit]
Description=Tikpal Speaker Chromium Kiosk
After=network.target tikpal-web.service
Wants=tikpal-web.service

[Service]
Type=simple
User=@SERVICE_USER@
Group=@SERVICE_GROUP@
WorkingDirectory=@APP_DIR@
EnvironmentFile=-@APP_DIR@/.env.kiosk
Environment=XAUTHORITY=/home/@SERVICE_USER@/.Xauthority
ExecStart=/usr/bin/env bash -lc 'APP_DIR="@APP_DIR@" /usr/bin/startx "@APP_DIR@/deploy/chromium/start-tikpal-kiosk-session.sh" -- :0 -br -nocursor'
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
